# chezmoi/private_dot_config/systemd/user/auto-fetch-repos.timer
[Unit]
Description=Automatically `git fetch` repos every 5 minutes

[Timer]
OnCalendar=*:5/10

[Install]
WantedBy=timers.target
