# home/private_dot_config/systemd/user/block-binary.timer
[Unit]
Description=Automatically kill binaries in blocked dirs

[Timer]
OnCalendar=*-*-* *:*:0/10

[Install]
WantedBy=timers.target
